# docker/ETL_python/crontab
# START CRON JOB
0 0 * * * python /etl.py 
# END CRON JOB