# docker/ETL_python/crontab
# START CRON JOB
#job a las 0am todo los dias
#0 0 * * * python /etl.py 

#job cada 15minutos
*/15 * * * * python /etl.py 

# END CRON JOB